# etc/init/reprepro_watch.conf
start on [2345]
stop on [06]

respawn

chdir /srv/reprepro_watch
exec supervisord -u ubuntu -c /srv/reprepro_watch/etc/supervisord.conf -n
